# deploy/systemd/essl-sync@.service
[Unit]
Description=eSSL realtime sync for device %i
After=network.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=iqac
Group=iqac
WorkingDirectory=/home/iqac/IDCS-Restart/backend
# If you store environment in backend/.env, uncomment the following line
# EnvironmentFile=/home/iqac/IDCS-Restart/backend/.env

# ExecStart uses wrapper; %i is instance token (e.g. 192.168.81.80_4370)
ExecStart=/usr/local/bin/essl-sync-wrapper %i
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
