# deploy/whatsapp.service
[Unit]
Description=WhatsApp OTP Service for IDCS
After=network.target

[Service]
Type=simple
User=iqac
Group=iqac
WorkingDirectory=/home/iqac/IDCS-Restart/whatsapp-server
ExecStart=/usr/bin/node /home/iqac/IDCS-Restart/whatsapp-server/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="NODE_ENV=production"
Environment="WA_API_KEY=IQAC_SECRET_123"
Environment="OBE_WHATSAPP_API_KEY=IQAC_SECRET_123"
Environment="CHROME_PATH=/usr/bin/google-chrome"

[Install]
WantedBy=multi-user.target
